# cacti/conf/cron.d
*/5 * * * * root php /opt/cacti/poller.php
